# extensions/kindletodo/kindletodo.upstart.conf
# Kindle Todo — Upstart service (kiosk).
# Install on the Kindle's ROOT filesystem at /etc/upstart/kindletodo.conf:
#
#   mntroot rw
#   cp kindletodo.upstart.conf /etc/upstart/kindletodo.conf
#   mntroot ro
#   initctl reload-configuration
#
# It starts once the reader framework is up, then boot-image.sh takes over the
# screen. `respawn` restarts the loop if it ever dies.

description "Kindle Todo - image mode kiosk"
author "dalager"

start on started framework
respawn
respawn limit 10 120

exec /bin/sh /mnt/us/extensions/kindletodo/bin/boot-image.sh
